# deploy/systemd/dm-watchdog-disk.service
[Unit]
Description=Dungeon Master disk-space watchdog (one-shot)
Documentation=file:///opt/dungeon-master/deploy/watchdogs/disk-space.sh

[Service]
Type=oneshot
User=dungeonmaster
Group=dungeonmaster
EnvironmentFile=-/etc/dungeon-master/env
ExecStart=/opt/dungeon-master/deploy/watchdogs/disk-space.sh
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/dungeon-master
