# deploy/systemd/dm-watchdog-flux.timer
[Unit]
Description=Run the FLUX-health watchdog every 5 minutes
Documentation=file:///opt/dungeon-master/deploy/watchdogs/flux-health.sh

[Timer]
# OnBootSec gives systemd time to settle (10s) before the first
# probe so a fresh boot doesn't fire on services still starting up.
# OnUnitActiveSec=5min keeps a steady cadence regardless of how
# long the watchdog itself took.
OnBootSec=30s
OnUnitActiveSec=5min
AccuracySec=10s
# Persist the last-run timestamp across reboots so the cadence
# survives systemd restarts.
Persistent=true

[Install]
WantedBy=timers.target
